# deploy/boot-progress-arm.service
# ==============================================================================
# File:    boot-progress-arm.service
# Purpose: Honest-instrument arm/reader.  Runs once at boot BEFORE the app:
#          reads the prior breadcrumb trail -> derives the clean/crash verdict
#          -> writes startup_log -> re-arms a fresh trail for this boot.
# Author:  (implementation plan 2026-05-17)
# Created: 2026-05-17 -- Sprint 38 T11 (honest boot-progress instrument)
#
# Why
# ---
# This unit is the SOLE startup_log writer (the legacy journal-scan canary was
# deleted in the T10 cutover).  It must run after the filesystem is available
# (After=local-fs.target) and strictly Before=eclipse-obd.service so the
# verdict for the prior boot is recorded and the trail re-armed before the app
# starts emitting this boot's milestones.
#
# US-277 silent-import guard (mandatory, not optional): WorkingDirectory AND
# Environment=PYTHONPATH=/home/mcornelison/Projects/Eclipse-01 MUST both be
# set.  `python -m src.pi.diagnostics.boot_progress` only resolves when the
# repo root is on PYTHONPATH; omitting it makes the import fail silently and
# the unit becomes a no-op -- the prior boot's verdict would never be written.
#
# Installation (manual parity reference for hand drills):
#   sudo cp deploy/boot-progress-arm.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable boot-progress-arm.service
# ==============================================================================

[Unit]
Description=Eclipse OBD boot-progress arm/reader (honest instrument)
After=local-fs.target
Before=eclipse-obd.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=mcornelison
WorkingDirectory=/home/mcornelison/Projects/Eclipse-01
Environment=PYTHONPATH=/home/mcornelison/Projects/Eclipse-01
ExecStart=/home/mcornelison/obd2-venv/bin/python -m src.pi.diagnostics.boot_progress --arm --file /home/mcornelison/Projects/Eclipse-01/data/boot_progress --db /home/mcornelison/Projects/Eclipse-01/data/obd.db --nas-dir /mnt/projects/O/OBD2v2/boot-progress --nas-enabled

[Install]
WantedBy=multi-user.target
