# deploy/eclipse-obd.service
# ==============================================================================
# File: eclipse-obd.service
# Purpose: systemd service file for Eclipse OBD-II Performance Monitor (Pi tier)
# Author: Ralph Agent
# Creation Date: 2026-01-23
# Copyright: (c) 2026 Michael Cornelison. All rights reserved.
#
# Installation (on the Pi, as root):
#   sudo cp eclipse-obd.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable eclipse-obd
#   sudo systemctl start eclipse-obd
#
# Or let install-service.sh do it for you — it substitutes paths/user/venv
# based on flags and matches the layout below.
#
# Defaults match the Sprint 10 Pi-crawl hostname migration:
#   User=mcornelison
#   WorkingDirectory=/home/mcornelison/Projects/Eclipse-01
#   Venv=/home/mcornelison/obd2-venv  (mirrors the server's ~/obd2-server-venv)
#   ExecStart=.../python src/pi/main.py
# ==============================================================================

[Unit]
Description=Eclipse OBD-II Performance Monitor
Documentation=https://github.com/your-repo/eclipse-obd
After=network.target bluetooth.target
Wants=bluetooth.target

# Start service after these are available
# network.target - network is configured
# bluetooth.target - Bluetooth subsystem is ready

# Flap protection — belongs in [Unit] (modern systemd warns if in [Service]).
# Allow up to StartLimitBurst restarts per StartLimitIntervalSec; beyond that,
# systemd gives up and the service stays failed until 'systemctl reset-failed'.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# Service type: simple means ExecStart is the main process
Type=simple

# User to run the service as (matches the Pi's interactive SSH user)
User=mcornelison

# Working directory — where src/pi/main.py lives after deploy-pi.sh sync
WorkingDirectory=/home/mcornelison/Projects/Eclipse-01

# Environment variables
# PATH includes the dedicated Pi venv (parallel to server's ~/obd2-server-venv)
Environment=PATH=/home/mcornelison/obd2-venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

# Main command to start the application
# src/pi/main.py resolves its default --config to config.json at repo root via
# Path(__file__), so no explicit --config flag is required here.
ExecStart=/home/mcornelison/obd2-venv/bin/python src/pi/main.py

# Restart policy
# on-failure: restart only if exit code is non-zero
Restart=on-failure

# Wait 10 seconds before restarting (allows Bluetooth to stabilize)
RestartSec=10

# (flap protection keys StartLimitIntervalSec / StartLimitBurst live in
#  [Unit] above — modern systemd ignores them here and logs a warning)

# Logging configuration
# StandardOutput and StandardError go to the systemd journal by default —
# view with: sudo journalctl -u eclipse-obd -f
# (The legacy StandardOutput=append directive was removed post-reorg;
#  journalctl is the single source of truth for Pi runtime logs.)

# Resource limits (optional, uncomment if needed)
# MemoryMax=512M
# CPUQuota=50%

# Security hardening (optional, uncomment if service is stable)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=read-only
# PrivateTmp=true

[Install]
# Start on multi-user.target (normal boot)
WantedBy=multi-user.target
