# deploy/orphan-cleanup.service
# ==============================================================================
# File:    orphan-cleanup.service
# Purpose: Nightly cleanup of NULL-drive_id realtime_data rows (US-322 / B-072)
# Author:  Rex (Ralph Agent)
# Created: 2026-05-11 -- Sprint 32 US-322 (V0.27.6)
#
# Why
# ---
# Spool 2026-05-11 audit found 61,293 NULL-drive_id rows accumulated in
# Pi-side realtime_data since Session 10's cleanup.  Sources are reconnect-
# loop polling, the I-019 DriveDetector warm-restart gap, and pre-DriveDetector
# grace-period rows.  Not poison (Spool's tuning analytics filter on
# drive_id IS NOT NULL), but unbounded growth is wasteful.  Approach 1
# (script + nightly timer) was chosen over Approach 2 (writer-side guard)
# for being simpler + reversible.
#
# Pairing
# -------
# This unit is paired with orphan-cleanup.timer which fires nightly at 03:00
# local.  Type=oneshot: each fire spins up a fresh Python interpreter, runs
# scripts/cleanup_orphan_realtime_data.py --execute, and exits.  The DB file
# is the only persistent state.
#
# Installation: auto-wired into deploy/deploy-pi.sh via
# step_install_orphan_cleanup_unit() (sync-if-changed install of .service +
# .timer + daemon-reload-on-change + enable --now).
#
# Verify:
#   systemctl list-timers orphan-cleanup.timer
#   journalctl -u orphan-cleanup.service -n 20
# ==============================================================================

[Unit]
Description=Eclipse OBD-II orphan realtime_data cleanup (US-322 / B-072)
Documentation=file:///home/mcornelison/Projects/Eclipse-01/scripts/cleanup_orphan_realtime_data.py
After=network.target

[Service]
Type=oneshot

# Run as the same user as the main eclipse-obd service so the SQLite file
# permissions are consistent (the live writer also runs as mcornelison).
User=mcornelison

WorkingDirectory=/home/mcornelison/Projects/Eclipse-01

Environment=PATH=/home/mcornelison/obd2-venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
# PYTHONPATH so the cleanup script (under scripts/) can import from src.* if
# a future refactor requires it.  The script is currently stdlib-only, but
# this matches the drain-forensics + eclipse-obd unit conventions.
Environment=PYTHONPATH=/home/mcornelison/Projects/Eclipse-01

ExecStart=/home/mcornelison/obd2-venv/bin/python scripts/cleanup_orphan_realtime_data.py \
    --db /home/mcornelison/Projects/Eclipse-01/data/obd.db \
    --age-hours 24 \
    --execute

# A failed cleanup MUST NOT page the operator; the next fire will retry
# and the live writer is unaffected.  The timer continues to drive
# subsequent fires regardless of any single-fire result.
SuccessExitStatus=0

# Tight resource caps so a runaway cleanup cannot starve the main service
# or the OBD poll loop running in a sibling cgroup.
MemoryMax=128M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
